# src/etc/systemd/system/ipsec.service.d/override.conf
[Unit]
ConditionPathExists=/run/swanctl/swanctl.conf
After=
After=vyos-router.service

[Service]
Environment="SWANCTL_DIR=/run/swanctl"
