# src/etc/systemd/system/kea-dhcp4-server.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
ExecStart=
ExecStart=/usr/sbin/kea-dhcp4 -c /run/kea/kea-dhcp4.conf
