# src/etc/systemd/system/ModemManager.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
ExecStart=
ExecStart=/usr/sbin/ModemManager --filter-policy=strict --log-level=INFO --log-timestamps --log-journal
