# src/systemd/vyos.target
[Unit]
Description=VyOS target
After=multi-user.target
