# src/etc/systemd/system/kea-dhcp6-server.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
ExecStart=
ExecStart=/usr/sbin/kea-dhcp6 -c /run/kea/kea-dhcp6.conf
