# src/etc/cron.d/check-wwan
*/5 * * * * root /usr/libexec/vyos/vyos-check-wwan.py
