# src/etc/systemd/system/kea-dhcp-ddns-server.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
ExecStart=
ExecStart=/usr/sbin/kea-dhcp-ddns -c /run/kea/kea-dhcp-ddns.conf
