# src/etc/systemd/system/LCDd.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
ExecStart=
ExecStart=/usr/sbin/LCDd -c /run/LCDd/LCDd.conf

