# src/etc/systemd/system/nginx.service.d/10-override.conf
[Unit]
After=
After=vyos-router.service
