# src/etc/systemd/system/serial-getty@.service.d/aftervyos.conf
[Service]
ExecStartPre=-/usr/libexec/vyos/init/vyos-config SERIAL
StandardOutput=journal+console
