# src/etc/systemd/system/salt-minion.service.d/override.conf
[Service]
SendSIGKILL=yes
FinalKillSignal=SIGQUIT
