# src/etc/update-motd.d/99-reboot
#!/bin/vbash
source /opt/vyatta/etc/functions/script-template
if [ -f /run/systemd/shutdown/scheduled ]; then
    echo
    run show reboot
fi
exit
