# src/etc/systemd/system/conntrackd.service.d/override.conf
[Unit]
After=
After=vyos-router.service
ConditionPathExists=/run/conntrackd/conntrackd.conf

[Service]
ExecStart=
ExecStart=/usr/sbin/conntrackd -C /run/conntrackd/conntrackd.conf
