# src/etc/systemd/system/conserver-server.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
Type=simple
ExecStart=
ExecStart=/usr/sbin/conserver -C /run/conserver/conserver.cf

