# src/etc/systemd/system/conserver-server.service.d/override.conf
[Unit]
After=
After=vyos-router.service
ConditionPathExists=/run/conserver/conserver.cf

[Service]
Type=simple
ExecStart=
ExecStart=/usr/sbin/conserver -M localhost -C /run/conserver/conserver.cf
Restart=on-failure
