# src/etc/systemd/system/hostapd@.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
WorkingDirectory=/run/LCDd
ExecStart=
ExecStart=/usr/sbin/LCDd -s 1 -f -c /run/LCDd/LCDd.conf

