# src/etc/systemd/system/atop.service.d/10-override.conf
[Service]
ExecStartPre=
ExecStart=
ExecStart=/bin/sh -c 'exec /usr/bin/atop ${LOGOPTS} -w "${LOGPATH}/atop.log" ${LOGINTERVAL}'
ExecStartPost=

